# omr-tracker/files/etc/uci-defaults/omr-tracker
#!/bin/sh
uci -q batch <<-EOF
    delete ucitrack.@omr-tracker[-1]
    add ucitrack omr-tracker
    set ucitrack.@omr-tracker[-1].init="omr-tracker"
    add_list ucitrack.@network[-1].affects="omr-tracker"
    commit ucitrack
EOF