# openmptcprouter/files/etc/hotplug.d/iface/00-nego
#!/bin/sh
# Needed for some device, like rtl8156 2.5G USB adapter
[ "$ACTION" = "ifup" ] && [ -n "$DEVICE" ] && ethtool -s $DEVICE autoneg on