# luci-theme-openmptcprouter/root/etc/uci-defaults/luci-theme-openmptcprouter
#!/bin/sh
uci batch <<-EOF
	set luci.themes.OpenMPTCProuter=/luci-static/openmptcprouter
	set luci.main.mediaurlbase=/luci-static/openmptcprouter
	commit luci
EOF
exit 0
