# openmptcprouter/files/etc/init.d/alwaysrw
#!/bin/sh /etc/rc.common

START=10

boot() {
	[ "$(mount | grep '/dev/root' | grep 'ext4' | grep '(ro')" != "" ] && mount -o remount,rw /dev/root
}
