# openmptcprouter/files/etc/init.d/fengshan
#!/bin/sh /etc/rc.common

START=99
start()
{
tem=$(cat /sys/class/thermal/thermal_zone0/temp)
echo $tem
if [ $tem -ge 60000 ] ; then
echo 1 > /sys/class/leds/green:fengshan/brightness;
echo "fengshandoup"
logger -t "fengshanup"
elif [ $tem -le 50000 ] ; then
echo 0 > /sys/class/leds/green:fengshan/brightness;
echo "fengshandown"
logger -t "$tem" "fengshandown"
fi
}