# openmptcprouter/files/etc/uci-defaults/1920-omr-network
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete network.wan1=interface
	set network.wan1=interface
	set network.wan1.ifname=eth0
	set network.wan1.proto=none
	set network.wan1.type=macvlan
	set network.wan1.ip4table=wan
	set network.wan1.multipath=master
	set network.wan1.defaultroute=0
	delete network.wan2=interface
	set network.wan2=interface
	set network.wan2.ifname=eth0
	set network.wan2.proto=none
	set network.wan2.type=macvlan
	set network.wan2.ip4table=wan
	set network.wan2.multipath=on
	set network.wan2.defaultroute=0
	commit network
EOF
rm -f /tmp/luci-indexcache
exit 0
