# openmptcprouter/files/etc/uci-defaults/1940-omr-dns
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	set unbound.@unbound[-1].listen_port=5353
	commit unbound
	add_list dhcp.@dnsmasq[-1].server="127.0.0.1#5353"
	commit dhcp
EOF
rm -f /tmp/luci-indexcache

exit 0
