# openmptcprouter/files/etc/uci-defaults/1980-omr-firewall
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	add firewall rule
	set firewall.@rule[-1].enabled='1'
	set firewall.@rule[-1].target='ACCEPT'
	set firewall.@rule[-1].name='Allow-All-Ping'
	set firewall.@rule[-1].proto='icmp'
	set firewall.@rule[-1].dest='*'
	set firewall.@rule[-1].src='*'
	set firewall.@rule[-1].icmp_type='echo-request'
	commit firewall
EOF
rm -f /tmp/luci-indexcache

exit 0
