# openmptcprouter/files/etc/uci-defaults/2000-omr-dhcp
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	set dhcp.lan.ra_default=0
	set dhcp.lan.dhcpv6=disabled
	set dhcp.lan.force=1
	commit dhcp
EOF
rm -f /tmp/luci-indexcache

exit 0
