# openmptcprouter/files/etc/uci-defaults/2050-omr-uhttpd
#!/bin/sh
uci -q batch <<-EOF >/dev/null
	set uhttpd.main.redirect_https=0
	set uhttpd.main.script_timeout=120
	set uhttpd.main.network_timeout=120
EOF
