# openmptcprouter/files/etc/uci-defaults/2060-omr-system
#!/bin/sh
uci -q batch <<-EOF >/dev/null
	set system.@system[-1].ttylogin=1
	commit system
	rpcd.@rpcd[0].timeout=60
	commit rpcd
EOF
exit 0