# openmptcprouter/files/etc/uci-defaults/2070-lcd4linux
#!/bin/sh 

/etc/init.d/lcd4linux stop >/dev/null 2>&1
/etc/init.d/lcd4linux disable >/dev/null 2>&1
exit 0
