# openmptcprouter/files/etc/uci-defaults/2080-omr-server
#!/bin/sh
if [ "$(uci -q get firewall.omr-server)" = "" ]; then
	uci -q batch <<-EOF >/dev/null
		set firewall.omr-server=include
		set firewall.omr-server.path=/etc/firewall.omr-server
		set firewall.omr-server.reload=1
		commit firewall
	EOF
fi

exit 0
