# qaa/qca/net/qca-hostap/files/qca-wpa-supplicant.init
#!/bin/sh /etc/rc.common
#
# Copyright (c) 2017 Qualcomm Technologies, Inc.
# All Rights Reserved.
# Confidential and Proprietary - Qualcomm Technologies, Inc.
#

#
# Copyright (c) 2016, The Linux Foundation. All rights reserved.
#

START=15

boot() {
	if [ -e "/var/run/wpa_supplicant-global.pid" ]
	then
		kill $(cat /var/run/wpa_supplicant-global.pid) &> /dev/null
		rm /var/run/wpa_supplicant-global.pid &> /dev/null
	fi
	wpa_supplicant -g /var/run/wpa_supplicantglobal -B -P /var/run/wpa_supplicant-global.pid
}

