# openmptcprouter/files/etc/crontabs/vnstat
* */12 * * * /etc/init.d/vnstat_backup backup