# luci-theme-lamatel/root/etc/uci-defaults/30_luci-theme-lametal
#!/bin/sh
uci batch <<-EOF
	set luci.themes.lamtel=/luci-static/lamatel
	set luci.main.mediaurlbase=/luci-static/lamatel
	commit luci
EOF
exit 0
