# omr-schedule/files/etc/uci-defaults/3000-omr-schedule
#!/bin/sh
if [ ! -f /etc/crontabs/root ] || [ "$(cat /etc/crontabs/root | grep omr-schedule)" = "" ]; then
	echo "*/5 * * * * /bin/omr-schedule" >> /etc/crontabs/root
fi
