# luci-theme-ezengreen/root/etc/uci-defaults/30_luci-theme-ezengreen
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Ezengreen=/luci-static/ezengreen
	set luci.main.mediaurlbase=/luci-static/ezengreen
	commit luci
EOF
exit 0
