# luci-theme-design/root/etc/uci-defaults/30_luci-theme-design
#!/bin/sh

uci get luci.themes.Design >/dev/null 2>&1 || \
uci batch <<-EOF
	set luci.themes.Design=/luci-static/design
	commit luci
EOF

exit 0
