# mptcp/files/etc/init.d/mptcp
#!/bin/sh /etc/rc.common

START=20
STOP=90

USE_PROCD=1

start_service() {
        procd_open_instance
        procd_set_param command hotplug-call net
        procd_set_param respawn
        procd_set_param watch network.globals
        procd_close_instance
}

reload_service() {
        ubus call network reload
}

service_running() {
        ubus -t 30 wait_for network
}
