# mptcp/files/etc/init.d/mptcp
#!/bin/sh /etc/rc.common

START=18

USE_PROCD=1

start_service() {
        /sbin/hotplug-call net
}

reload_service() {
        start
}

service_triggers() {
        procd_add_reload_trigger network
}
